# u2f-nfchidemu.service
[Unit]
Description=U2F-NFC HID Emulation

[Service]
ExecStart=/usr/bin/u2f-nfchidemu
Type=simple

